# config/systemd/cobalt-welcome.service
[Unit]
Description=CobaltOS Welcome App (first login)
Documentation=https://github.com/0xnullsect0r/CobaltOS
ConditionPathExists=!/var/lib/cobalt/welcome-done
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/cobalt-welcome
# After the wizard exits: mark done, then start the Flatpak install queue
ExecStopPost=/bin/sh -c 'mkdir -p /var/lib/cobalt && touch /var/lib/cobalt/welcome-done'
ExecStopPost=systemctl --user start cobalt-flatpak-install.service
Restart=no
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical-session.target
